# rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user start-up script.
#

echo "Setting backlight brightness"
echo 4 > /sys/class/backlight/acpi_video0/brightness
